# confs/uwsgi.conf
# description "uWSGI server for {{ PROJECT_NAME }}"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env DEPLOYMENT_TARGET={{ DEPLOYMENT_TARGET }}
env TUMBLR_CONSUMER_KEY={{ TUMBLR_APP_KEY }}
env TUMBLR_OAUTH_TOKEN={{ TUMBLR_OAUTH_TOKEN }}
env TUMBLR_OAUTH_TOKEN_SECRET={{ TUMBLR_OAUTH_TOKEN_SECRET }}
env TUMBLR_APP_SECRET={{ TUMBLR_APP_SECRET }}
env AWS_SECRET_ACCESS_KEY={{ AWS_SECRET_ACCESS_KEY }}
env AWS_ACCESS_KEY_ID={{ AWS_ACCESS_KEY_ID }}

script
    /usr/local/bin/uwsgi --ini /home/ubuntu/apps/{{ REPOSITORY_NAME }}/repository/{{ PROJECT_SLUG }}.app.ini
end script
