# confs/uwsgi.conf
# description "uWSGI server for {{ PROJECT_SLUG }}"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
    . /etc/environment
    /usr/local/bin/uwsgi --ini {{ SERVER_REPOSITORY_PATH }}/{{ PROJECT_FILENAME }}.app.ini
end script
