# crontab
* * * * * ubuntu date >> /tmp/cron-test.txt 
