# crontab
* * * * * bash /home/ubuntu/apps/app-template/repository/cron.sh fab $DEPLOYMENT_TARGET cron_test >/var/log/app-template-crontab.log 2>&1