# crontab
* * * * * ubuntu bash /home/ubuntu/apps/app-template/repository/cron.sh fab $DEPLOYMENT_TARGET cron_test 
