# crontab
* * * * * ubuntu /bin/bash /home/ubuntu/apps/app-template/repository/run_on_server.sh fab $DEPLOYMENT_TARGET cron_test >/var/log/app-template-crontab.log 2>&1
