# crontab
* * * * * ubuntu echo $DEPLOYMENT_TARGET > /tmp/cron-test.txt 
