# confs/rendered/app-template.uwsgi.conf
description "uWSGI server for App Template"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env DEPLOYMENT_TARGET=staging
env TUMBLR_OAUTH_TOKEN=H8ISk2DIcMR0Bq2dsHCCcDMY6XpKAWDOApDQWCGWKLJGMDc2dw
env TUMBLR_OAUTH_TOKEN_SECRET=6Jz1pI06cqum1MoohBEgV4g38Z8yVYOaGrXP4B3qp0xukeo93B
env TUMBLR_APP_SECRET=QYQ6xuMMYzRmovnkiN1t5V0pLoeTPTYzNrMt1WH3gLDu3cm7XA
env AWS_SECRET_ACCESS_KEY=4lzfhbCG8JLzu4Ynm9ijhUQ4aDXuwkNwb8ZVylGc
env AWS_ACCESS_KEY_ID=AKIAIMHRQT3YCGQBWKAQ

script
    /usr/local/bin/uwsgi \
    --virtualenv /home/ubuntu/apps/app-template/virtualenv \
    --chdir /home/ubuntu/apps/app-template/repository \
    --file public_app.py \
    --chmod-socket 777 \
    --chown-socket www-data:www-data \
    --catch-exceptions \
    --touch-reload /home/ubuntu/apps/app-template/repository/public_app.py \
    --callable app \
    --logto /var/log/uwsgi.log \
    --die-on-term \
    -p 15 \
    -M \
    -s /tmp/app-template.sock
end script